# edge_camera_node/systemd/traffic-camera-node.service
[Unit]
Description=Traffic Warning Edge Camera Node
# Do not block startup on network-online; node can run and recover when Wi-Fi appears later.
After=network.target NetworkManager.service avahi-daemon.service
Wants=NetworkManager.service avahi-daemon.service
StartLimitIntervalSec=60
StartLimitBurst=10

[Service]
Type=simple
User=__RUN_USER__
WorkingDirectory=__WORKDIR__
ExecStart=__PYTHON__ -m traffic_camera_node.main --config __CONFIG__
Restart=on-failure
RestartSec=3
# Lower scheduling priority slightly so networking/user login stays responsive under load.
Nice=5
IOSchedulingClass=best-effort
IOSchedulingPriority=6
TasksMax=512
KillSignal=SIGTERM
TimeoutStopSec=20
Environment=PYTHONUNBUFFERED=1

[Install]
WantedBy=multi-user.target
